# systemd/system/key-locker.path
[Unit]
Description=Run key-locker recv

[Path]
PathModified=/tmp/key-locker.path

[Install]
WantedBy=multi-user.target