# systemd/system/key-locker.service
[Unit]
Description=Run key-locker recv

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/key-locker recv